# deployment/backend/systemd/careersignals.service
[Unit]
Description=CareerSignals production Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/careersignals/current
ExecStart=/usr/bin/docker compose --env-file /opt/careersignals/current/release.env -f /opt/careersignals/current/docker-compose.production.yml up -d
ExecReload=/usr/bin/docker compose --env-file /opt/careersignals/current/release.env -f /opt/careersignals/current/docker-compose.production.yml up -d
ExecStop=/usr/bin/docker compose --env-file /opt/careersignals/current/release.env -f /opt/careersignals/current/docker-compose.production.yml stop
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
